# templates/etc/systemd/system/minecraft@.socket
[Unit]
PartOf=minecraft@.service

[Socket]
ListenFIFO=%t/minecraft/%i.stdin